# deploy/ansible/roles/hort_systemd/templates/hort-server.service.j2
# hort-server.service — managed by Ansible (hort-systemd role).
# The hort artifact registry API server.  Runs as non-root user `hort`
# with the isolation directives below (replacing the container boundary
# that the Podman flavor uses).

[Unit]
Description=hort-server
Requires=hort-migrate.service postgresql.service
After=hort-migrate.service postgresql.service

[Service]
User={{ hort_user }}
Group={{ hort_user }}
ExecStart={{ hort_systemd_server_bin }} serve
EnvironmentFile={{ hort_systemd_env_dir }}/hort-server.env

# Isolation boundary (replaces the container):
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths={{ hort_systemd_cas_dir }}
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# ── M1 additional hardening (security audit 2026-06) ──────────────────────────
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
PrivateDevices=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0077
LimitNOFILE=65536
TasksMax=4096
# NOTE: MemoryDenyWriteExecute is intentionally OMITTED.
# wasmtime uses a JIT that requires W+X memory mappings at runtime; MDWX
# would cause an immediate SIGSYS on the first WASM module load.  Do NOT
# add MemoryDenyWriteExecute to this unit without first replacing wasmtime
# with an AOT-only runtime that makes no W+X mappings.

Restart=always

[Install]
WantedBy=multi-user.target
