# deploy/ansible/roles/hort_systemd/templates/hort-worker.service.j2
# hort-worker.service — managed by Ansible (hort-systemd role).
# The hort background worker: scan dispatch, quarantine-release sweep,
# prefetch tick, etc.  No listening ports; trivy + osv-scanner must be
# on PATH (installed by hort-systemd role Step 4).

[Unit]
Description=hort-worker
Requires=hort-migrate.service postgresql.service
After=hort-migrate.service postgresql.service

[Service]
User={{ hort_user }}
Group={{ hort_user }}
ExecStart={{ hort_systemd_worker_bin }}
EnvironmentFile={{ hort_systemd_env_dir }}/hort-worker.env

# Isolation boundary (replaces the container):
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths={{ hort_systemd_cas_dir }}
CapabilityBoundingSet=
# Worker needs AF_UNIX for the local postgres socket and AF_INET/AF_INET6
# for scanner advisory API calls and upstream registry fetches.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
Restart=always

[Install]
WantedBy=multi-user.target
